# webhook_deploy.service
[Unit]
Description=GitHub Webhook Auto Deploy Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/volume2/docker/xiaole-ai
Environment="WEBHOOK_SECRET=your-webhook-secret-change-me"
Environment="DB_USER=your_db_user"
Environment="DB_PASS=your_db_password"
Environment="DEEPSEEK_API_KEY=your_deepseek_api_key"
Environment="QWEN_API_KEY=your_qwen_api_key"
Environment="BAIDU_APP_ID=your_baidu_app_id"
Environment="BAIDU_API_KEY=your_baidu_api_key"
Environment="BAIDU_SECRET_KEY=your_baidu_secret_key"
ExecStart=/usr/bin/python3 /volume2/docker/xiaole-ai/webhook_deploy.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
